# usr/local/etc/rc.local
#!/bin/sh -e

# Run S*.rc sequence in /usr/local/etc/rc.local.d
for i in $(find /usr/local/etc/rc.local.d -type f -name "S*" | sort); do
	. ${i}
done

# Send info to dmesg
echo 'System booted' > /dev/kmsg

exit 0
